# data/niri-screenshare.service
[Unit]
Description=Portal service (niri backend)
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.impl.portal.desktop.niri
ExecStart=/usr/lib/niri-screenshare
Restart=on-failure

[Install]
WantedBy=graphical-session.target
